# deploy/systemd/salesman-daily.service
[Unit]
Description=PlausiDen-Salesman daily orchestration runbook (run by .timer)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=salesman
Group=salesman
EnvironmentFile=/etc/salesman.env
# Required by daily.sh — operator sets these in /etc/salesman.env or
# overrides by editing the unit. Without them the script bails fast.
Environment=SALESMAN_AUTOPILOT=1
WorkingDirectory=/opt/salesman/build/PlausiDen-Salesman
ExecStart=/opt/salesman/build/PlausiDen-Salesman/scripts/salesman-daily.sh
TimeoutStartSec=30min
StandardOutput=append:/opt/salesman/logs/daily.log
StandardError=append:/opt/salesman/logs/daily.log

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/opt/salesman/data /opt/salesman/logs
