# deploy/systemd/salesman-inbox-poll.service
[Unit]
Description=PlausiDen-Salesman IMAP inbox poll (run by .timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=salesman
Group=salesman
EnvironmentFile=/etc/salesman.env
ExecStart=/opt/salesman/bin/salesman inbox-poll

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/opt/salesman/data /opt/salesman/logs
